# systemd/printerfacts.service.in
[Unit]
Description=The Printerfacts service

[Service]
DynamicUser=yes
Environment=PORT=32042
Environment=RUST_LOG=info
WorkingDir=@printerfacts@
ExecStart=@printerfacts@/bin/printerfacts

[Install]
WantedBy=multi-user.target
